# rbackup.service
[Unit]
Description=Redis Backup service

[Service]
Type=simple
ExecStart=/opt/redis_backup/rbackup.py

[Install]
WantedBy=multi-user.target
